# media-libs/display-link-bin/files/amd64/90-smiusbdisplay.rules
# Copyright (c) 2018 - 2020 SiliconMotion Inc.
# File autogenerated by udev-installer.sh script

ACTION=="add", SUBSYSTEM=="usb",\
 DRIVERS=="usb", ATTR{idVendor}=="090c",\
 MODE="0660", GROUP="video", \
 IMPORT{builtin}="usb_id", ENV{SMIUSBDISPLAY_DEVNAME}="$env{DEVNAME}",\
 ENV{SMIUSBDISPLAY_DEVICE_ID}="$env{ID_BUS}-$env{BUSNUM}-$env{DEVNUM}-$env{ID_SERIAL}",\
 SYMLINK+="siliconmotion/by-id/$env{SMIUSBDISPLAY_DEVICE_ID}"

ACTION=="add", SUBSYSTEM=="usb",\
 DRIVERS=="usb", ATTR{idVendor}=="090c",\
 ATTR{bDeviceClass}=="ef" , IMPORT{parent}="SMIUSBDISPLAY*",\
 RUN+="/opt/siliconmotion/smi-udev.sh add $root $env{SMIUSBDISPLAY_DEVICE_ID} $env{SMIUSBDISPLAY_DEVNAME}"

ACTION=="remove",\
 ENV{PRODUCT}=="90c/*", RUN+="/opt/siliconmotion/smi-udev.sh remove $root $env{DEVNAME}"
